# systemd/lq45-signal-bot.service
[Unit]
Description=LQ45 AI Signal Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/lq45-signal-bot
Environment=PATH=/root/lq45-signal-bot/venv/bin
ExecStart=/root/lq45-signal-bot/venv/bin/python3 /root/lq45-signal-bot/main.py
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
